# Docker/conf/crontabs/www-data
*/5  *  *  *  * php -f /var/www/nextcloud/releases/current/cron.php
